# airpods-voice.service
[Unit]
Description=AirPods Voice Control
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/airpods-button-listener
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
